# docker/config/etc/supervisord/conf.d/api/queue-work.conf
[program:queue-work]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /app/artisan queue:work --queue=create-report-queue
autostart=true
autorestart=true
user=root
numprocs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=7200
