# docker/config/etc/supervisord/conf.d/api/schedule.conf
[program:cron]
command=/usr/sbin/crond -f -L /tmp/cron.log
user=root
autostart=true
autorestart=true
